# deploy/baxter.service
# Baxter Burgundy -- bring the docker-compose fleet up on boot and give a single
# `systemctl` handle for the stack. One-time install steps: see deploy/README.md.
#
# This is a ONESHOT unit, not a long-running process: `make run-gmail` starts the
# containers detached (compose `up -d`) and returns, so the unit sits in
# "active (exited)". Per-container liveness -- a crash, an OOM kill, a Docker
# daemon restart, or a reboot -- is handled by each service's
# `restart: unless-stopped` policy in compose.yaml plus an enabled Docker daemon.
# systemd's only jobs here are "bring the stack up at boot" and "stop it
# gracefully on `systemctl stop baxter`".

[Unit]
Description=Baxter Burgundy agent fleet (docker compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/baxter

# Type=oneshot already runs with NO start timeout by default (systemd.service(5)),
# which is what we want: the first boot after a fresh deploy builds the app +
# codapi images and can take several minutes. Stated explicitly so it survives a
# future switch away from oneshot. Set a finite value instead if you'd rather a
# genuinely wedged build fail than hang boot. (Later boots hit the layer cache.)
TimeoutStartSec=infinity

# PROJECT is passed as a make ARGUMENT, never an env prefix: the Makefile's
# `PROJECT := $(notdir $(CURDIR))` ignores an env var, so an argument is the only
# way to pin the fleet name. (With WorkingDirectory=/opt/baxter it would also
# default to "baxter" from the dir name, but pinning it keeps the unit correct
# regardless of where the repo is checked out.)
ExecStart=/usr/bin/make run-gmail PROJECT=baxter
ExecStop=/usr/bin/make stop PROJECT=baxter

# The user that owns /opt/baxter and is a member of the `docker` group, so make
# can drive the Docker socket. Set it AFTER install via a drop-in, NOT here:
# `sudo systemctl edit baxter` -> [Service] / User=<box user>. Editing this
# tracked file in place would dirty the working tree and block `make deploy`
# (the clean-tree guard). Left as CHANGEME on purpose: a missing override fails
# loudly (systemd: "no such user") instead of silently running as root.
User=CHANGEME

[Install]
WantedBy=multi-user.target
